# portal/scripts/ipsum-refresh.timer
[Unit]
Description=Daily refresh of IPsum IP threat intelligence feed
Requires=ipsum-refresh.service

[Timer]
# Run at 04:30 UTC daily (feed is typically updated ~03:00 UTC)
OnCalendar=*-*-* 04:30:00 UTC
# Run immediately if the last scheduled run was missed (e.g. host was off)
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
